# piculator-localhost_1.0-1/DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable piculator-localhost
systemctl service piculator-localhost start