# piculator-localhost_1.0-1/DEBIAN/prerm
#!/bin/bash
systemctl stop piculator-localhost
systemctl disable piculator-localhost